# shell/csu-autoauth.service.example
[Unit]
Description=CSU Auto Auth

[Service]
ExecStart=%h/.local/bin/csu-autoauth.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target